# bun/systemd/katana.service
[Unit]
Description=Katana Module Management Server
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/katana serve --tls --host 0.0.0.0
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/katana /etc/nginx/sites-available /etc/nginx/sites-enabled

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
